# examples/ocp.example.app/files/publicurl/fixpublicurl.service
[Unit]
Description=Fix OpenShift publicURL
After=atomic-openshift-master.service

[Service]
ExecStart=/usr/local/bin/update_publicURL.sh
Type=simple

[Install]
WantedBy=default.target
